# daily_report.timer
[Unit]
Description=Run NeiroTolikBot daily analytics report

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
