# daily_report.service
[Unit]
Description=NeiroTolikBot daily analytics report
Wants=daily_report.timer

[Service]
Type=oneshot
WorkingDirectory=/root/tolikNeiroBot/NeiroTolikBot
EnvironmentFile=/root/tolikNeiroBot/NeiroTolikBot/.env
ExecStart=/root/tolikNeiroBot/NeiroTolikBot/venv/bin/python /root/tolikNeiroBot/NeiroTolikBot/scripts/daily_report.py
